# etc/systemd/system/docker-prune.timer
[Unit]
Description=Weekly Docker disk reclamation (docker system prune)

[Timer]
OnCalendar=*-*-* 05:00:00
# Daily so the until=168h floor collects week-old junk within a day of it
# becoming eligible, rather than letting it linger up to another week. This host
# (Crostini) is frequently suspended; without Persistent a missed window is
# lost, so Persistent runs the catch-up on next wake. 05:00 keeps it an hour
# clear of the dev-cache-gc sweep on the days they overlap.
Persistent=true

[Install]
WantedBy=timers.target
